# etc/supervisor.conf
[program:bundy]
command=bash -c "/home/pv/bin/bundy"
process_name=%(program_name)s
startretries=10
stopsignal=QUIT
stopasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
